# files/.config/profile.d/source.sh
source $HOME/.export
source $HOME/.langs
source $HOME/.aliases
