# files/.config/profile.d/work.sh
WORK_FILE=$HOME/Life/work.sh
if [ -f $WORK_FILE ];
then
    source $WORK_FILE
fi
